# deployment/linux/ml-gestao-api.service
[Unit]
Description=ML Gestao Total API
After=network.target

[Service]
Type=simple
User=mlgestao
Group=mlgestao
WorkingDirectory=/opt/ml-gestao-total
Environment=NODE_ENV=production
Environment=PORT=3001
EnvironmentFile=/opt/ml-gestao-total/.env.local
ExecStart=/usr/bin/node /opt/ml-gestao-total/server-dist/server/core/index.js
Restart=on-failure
RestartSec=5
TimeoutStopSec=20
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
UMask=0077
ReadWritePaths=/opt/ml-gestao-total/data /opt/ml-gestao-total/backups

[Install]
WantedBy=multi-user.target
